# templates/conntrack-logging.service.j2
[Unit]
Description=Conntrack table logging service

[Service]
ExecStart={{ conntrack_logging_conntrack_path }} {{ conntrack_logging_conntrack_options }}
Restart=always
